# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-ca-add-package/run
#!/usr/bin/with-contenv bash

cp /data/ca/* /usr/local/share/ca-certificates/ -f
update-ca-certificates
